# system_files/shared/usr/share/kyanite/quadlets/ollama-vulkan.container
[Unit]
Description=Ollama (containerized, Vulkan)
After=network-online.target
Wants=network-online.target
Conflicts=ollama.service ollama-rocm.service

[Container]
Image=docker.io/ollama/ollama:latest
ContainerName=ollama-vulkan
AutoUpdate=registry
Network=ai
NetworkAlias=ollama
PublishPort=127.0.0.1:11434:11434
Volume=ollama-data:/root/.ollama
AddDevice=/dev/dri
GroupAdd=keep-groups
Environment=OLLAMA_VULKAN=true

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target
