# system_files/shared/usr/share/kyanite/quadlets/open-webui.container
[Unit]
Description=Open WebUI (containerized)
After=network-online.target
Wants=network-online.target

[Container]
Image=ghcr.io/open-webui/open-webui:main
ContainerName=open-webui
AutoUpdate=registry
Network=ai
PublishPort=127.0.0.1:8080:8080
Volume=open-webui-data:/app/backend/data
Environment=WEBUI_HOST=0.0.0.0
Environment=OLLAMA_BASE_URL=http://ollama:11434

[Service]
Restart=on-failure
TimeoutStartSec=900

[Install]
WantedBy=default.target
